# build/packageroot/DEBIAN/postinst
#!/bin/sh

systemctl enable unifi-blueberry-core
systemctl start unifi-blueberry-core
