# build/packageroot/DEBIAN/prerm
#!/bin/sh

systemctl disable unifi-blueberry-core
systemctl stop unifi-blueberry-core
