# build/packageroot/lib/systemd/system/unifi-blueberry-core.service
[Unit]
Description=UniFi Blueberry core service
After=network.target

[Service]
ExecStart=/usr/bin/unifi-blueberry-core

[Install]
WantedBy=multi-user.target
